# battstat.service
[Unit]
Description=BattStat — UPS Battery Health Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=battstat
Group=battstat
WorkingDirectory=/opt/battstat
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=battstat

Environment=NODE_ENV=production
Environment=PORT=3000
# Listens on all interfaces — restrict access with ufw/firewalld on the host.
# Set to 127.0.0.1 if you are putting nginx/caddy in front as a reverse proxy.
Environment=HOST=0.0.0.0
# Uncomment and set to your public URL if using a reverse proxy:
# Environment=ALLOWED_ORIGIN=https://battstat.internal.company.com
# Set HTTPS=true when serving behind a TLS-terminating reverse proxy.
# This enables the Secure flag on session cookies. Do NOT set over plain HTTP.
# Environment=HTTPS=true
Environment=DB_PATH=/opt/battstat/data/battstat.db

# ── Systemd hardening ─────────────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/battstat/data
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
# NOTE: MemoryDenyWriteExecute and SystemCallFilter are intentionally omitted.
# Node.js/V8 requires W+X memory for JIT compilation — those flags cause a
# SIGTRAP crash on modern kernels. All other hardening above remains effective.

[Install]
WantedBy=multi-user.target
